# packages/botd/botd-7.tar.gz/botd-7/init.d/botd.service
[Unit]
Description=botd

[Service]
ExecStart=/usr/local/bin/botd -w /var/lib/botd -a /var/log/botd -z -l debug

[Install]
WantedBy=multi-user.target
