# packages/botd/botd-5.tar.gz/botd-5/init.d/botd.service
[Unit]
Description=botd

[Service]
ExecStart=/usr/local/bin/botd -d /var/lib/botd -a /var/log/botd -l debug

[Install]
WantedBy=multi-user.target
